# tools/freebsd/etc.rc.d/clid.in
#!/bin/sh
# PROVIDE: clid
# REQUIRE: LOGIN FILESYSTEMS
# KEYWORD: shutdown

. /etc/rc.subr

# Creating the user:
#   pw useradd clid -d /nonexistent -g yans -s /usr/sbin/nologin
#
# Add the yansgroup if not already present
#   pw groupadd yans
#
#   pw -R <rootdir> ... for images/buildroots/jails
#
# The daemon drops user privs itself, which is why clid_user, clid_group
# is not used
: ${clid_enable:="no"}
: ${clid_duser:="clid"}
: ${clid_dgroup:="yans"}
: ${ystatedir:="@LOCALSTATEDIR@/yans"}
: ${clid_basepath:="${ystatedir}/clid"}

name="clid"
rcvar="clid_enable"
command="@BINDIR@/${name}"
command_args="-u '${clid_duser}' -g '${clid_dgroup}' -b '${clid_basepath}'"
start_precmd="clid_precmd"
pidfile="${clid_basepath}/clid.pid"

clid_precmd()
{
  install -m 775 -o root -g yans -d "${ystatedir}" || {
    err 1 "Failed to create ${ystatedir}"
  }

  install -m 770 -o "${clid_duser}" -g "${clid_dgroup}" \
      -d "${clid_basepath}" || {
    err 1 "Failed to create ${clid_basepath}"
  }

  return 0
}

load_rc_config $name
run_rc_command "$1"
