# tools/freebsd/etc.rc.d/clid.in
#!/bin/sh
# PROVIDE: clid
# REQUIRE: LOGIN FILESYSTEMS
# KEYWORD: shutdown

. /etc/rc.subr

: ${clid_enable:="no"}
: ${clid_duser:="clid"}
: ${clid_dgroup:="yans"}
: ${ystatedir:="@LOCALSTATEDIR@/yans"}
: ${clid_basepath:="${ystatedir}/clid"}

name="clid"
rcvar="clid_enable"
command="@BINDIR@/${name}"
command_args="-u '${clid_duser}' -g '${clid_dgroup}' -b '${clid_basepath}'"
start_precmd="clid_precmd"
pidfile="${clid_basepath}/clid.pid"

clid_precmd()
{
  install -m 775 -o root -g yans -d "${ystatedir}" || {
    err 1 "Failed to create ${ystatedir}"
  }

  install -m 770 -o "${clid_duser}" -g "${clid_dgroup}" \
      -d "${clid_basepath}" || {
    err 1 "Failed to create ${clid_basepath}"
  }

  if [ -f "$pidfile" ]; then
    local pid=$(check_pidfile "$pidfile" "$command")
    if [ -z "$pid" ]; then
      warn "Removing stale pidfile"
      rm "$pidfile"
    else
      # Shouldn't happen due to the way rc works, but check for it anyway
      err 1 "Already running at $pid"
    fi
  fi

  return 0
}

load_rc_config $name
run_rc_command "$1"
