# tools/freebsd/etc.rc.d/ethd.in
#!/bin/sh
# PROVIDE: ethd
# REQUIRE: LOGIN FILESYSTEMS
# KEYWORD: shutdown

. /etc/rc.subr

: ${ethd_enable:="no"}
: ${ethd_duser:="ethd"}
: ${ethd_dgroup:="yans"}
: ${ystatedir:="@LOCALSTATEDIR@/yans"}
: ${ethd_basepath:="${ystatedir}/ethd"}
: ${ethd_devfs_ruleset:=""}

name="ethd"
rcvar="ethd_enable"
command="@BINDIR@/${name}"
command_args="-u '${ethd_duser}' -g '${ethd_dgroup}' -b '${ethd_basepath}'"
start_precmd="ethd_precmd"
pidfile="${ethd_basepath}/ethd.pid"

mountpoint_is_type() {
  local _ mountpoint fstype
  while read -r _ mountpoint fstype _; do
    if [ "$1" = "$mountpoint" -a "$2" = "$fstype" ]; then
      return 0
    fi
  done <<EOF
  $(mount -p)
EOF
  return 1
}

ethd_precmd()
{
  [ -z "${ethd_devfs_ruleset}" ] && err 1 "No devfs ruleset for ethd"

  install -m 775 -o root -g yans -d "${ystatedir}" || {
    err 1 "Failed to create ${ystatedir}"
  }

  install -m 750 -o "${ethd_duser}" -g "${ethd_dgroup}" \
      -d "${ethd_basepath}" || {
    err 1 "Failed to create ${clid_basepath}"
  }

  # create the devfs mountpoint for bpf device
  install -m 750 -o "${ethd_duser}" -g "${ethd_dgroup}" \
      -d "${ethd_basepath}/dev" || {
    err 1 "Failed to create ${clid_basepath}/dev"
  }

  # mount devfs if not already mounted
  mountpoint_is_type "${ethd_basepath}/dev" devfs || {
    mount -t devfs -o "ruleset=${ethd_devfs_ruleset}" devfs \
        "${ethd_basepath}/dev" || {
      err 1 "Failed to mount ethd devfs"
    }

    # sanity check devfs rule: nonexistent rulesets means mounting with the
    # default one
    [ -e "${ethd_basepath}/dev/mem" -o -e "${ethd_basepath}/dev/kmem" ] &&
      err 1 "a too permissive devfs is mounted for ethd"
  }

  if [ -f "$pidfile" ]; then
    local pid=$(check_pidfile "$pidfile" "$command")
    if [ -z "$pid" ]; then
      warn "Removing stale pidfile"
      rm "$pidfile"
    else
      # Shouldn't happen due to the way rc works, but check for it anyway
      err 1 "Already running at $pid"
    fi
  fi

  return 0
}

load_rc_config $name
run_rc_command "$1"
