# tools/freebsd/etc.rc.d/knegd.in
#!/bin/sh
# PROVIDE: knegd
# REQUIRE: LOGIN FILESYSTEMS
# KEYWORD: shutdown

. /etc/rc.subr

# Creating the user:
#   pw useradd knegd -d /nonexistent -g yans -s /usr/sbin/nologin
#
# Add the yansgroup if not already present
#   pw groupadd yans
#
#   pw -R <rootdir> ... for images/buildroots/jails
#
# The daemon drops user privs itself, which is why knegd_user, knegd_group
# is not used
: ${knegd_enable:="no"}
: ${knegd_duser:="knegd"}
: ${knegd_dgroup:="yans"}
: ${ystatedir:="@LOCALSTATEDIR@/yans"}
: ${knegd_basepath:="${ystatedir}/knegd"}

name="knegd"
rcvar="knegd_enable"
command="@BINDIR@/${name}"
command_args="-u '${knegd_duser}' -g '${knegd_dgroup}' -b '${knegd_basepath}'"
start_precmd="knegd_precmd"
pidfile="${knegd_basepath}/knegd.pid"

knegd_precmd()
{
  install -m 775 -o root -g yans -d "${ystatedir}" || {
    err 1 "Failed to create ${ystatedir}"
  }

  install -m 770 -o "${knegd_duser}" -g "${knegd_dgroup}" \
      -d "${knegd_basepath}" || {
    err 1 "Failed to create ${knegd_basepath}"
  }

  return 0
}

load_rc_config $name
run_rc_command "$1"
