# tools/freebsd/etc.rc.d/sc2.in
#!/bin/sh
# PROVIDE: sc2
# REQUIRE: LOGIN FILESYSTEMS
# KEYWORD: shutdown

. /etc/rc.subr

# Creating the user:
#   pw useradd sc2 -d @LOCALSTATEDIR@/sc2 -g yans -s /usr/sbin/nologin
#
# Add the yansgroup if not already present
#   pw groupadd yans
#
#   pw -R <rootdir> ... for images/buildroots/jails
#
# The daemon drops user privs itself, which is why sc2_user, sc2_group
# is not used
: ${sc2_enable:="no"}
: ${sc2_driver:="@BINDIR@/yansapi"}
: ${sc2_basepath:="@LOCALSTATEDIR@/sc2"}
: ${sc2_duser:="sc2"}
: ${sc2_dgroup:="yans"}

name="sc2"
rcvar="sc2_enable"
command="@BINDIR@/sc2"
command_args="-u ${sc2_duser} -g ${sc2_dgroup} -b ${sc2_basepath} ${sc2_driver}"
start_precmd="sc2_precmd"
pidfile="${sc2_basepath}/sc2.pid"

sc2_precmd()
{
  if [ ! -d "${sc2_basepath}" ]; then
    mkdir -p "${sc2_basepath}" || {
      warn "Failed to create sc2 directory"
      return 1
    }
  fi

  chmod 770 "${sc2_basepath}" || {
    warn "Failed to set correct permissions of ${sc2_basepath}"
    return 1
  }

  chown "${sc2_duser}:${sc2_dgroup}" "${sc2_basepath}" || {
    warn "Failed to set correct ownership of ${sc2_basepath}"
    return 1
  }

  return 0
}

load_rc_config $name
run_rc_command "$1"
