# tools/freebsd/etc.rc.d/stored.in
#!/bin/sh
# PROVIDE: stored
# REQUIRE: LOGIN FILESYSTEMS
# KEYWORD: shutdown

. /etc/rc.subr

# Creating the user:
#   pw useradd stored -d /nonexistent -g yans -s /usr/sbin/nologin
#
# Add the yansgroup if not already present
#   pw groupadd yans
#
#   pw -R <rootdir> ... for images/buildroots/jails
#
# The daemon drops user privs itself, which is why stored_user, stored_group
# is not used
: ${stored_enable:="no"}
: ${stored_duser:="stored"}
: ${stored_dgroup:="yans"}
: ${ystatedir:="@LOCALSTATEDIR@/yans"}
: ${stored_basepath:="${ystatedir}/stored"}

name="stored"
rcvar="stored_enable"
command="@BINDIR@/${name}"
command_args="-u ${stored_duser} -g ${stored_dgroup} -b '${stored_basepath}'"
start_precmd="stored_precmd"
pidfile="${stored_basepath}/stored.pid"

stored_precmd()
{
  install -m 775 -o root -g yans -d "${ystatedir}" || {
    err 1 "Failed to create ${ystatedir}"
  }

  install -m 770 -o "${stored_duser}" -g "${stored_dgroup}" \
      -d "${stored_basepath}" || {
    err 1 "Failed to create ${stored_basepath}"
  }

  return 0
}

load_rc_config $name
run_rc_command "$1"
