# tools/freebsd/etc.rc.d/yansapi.in
#!/bin/sh
# PROVIDE: yansapi
# REQUIRE: LOGIN FILESYSTEMS
# KEYWORD: shutdown

. /etc/rc.subr

: ${yansapi_enable:="no"}
: ${yansapi_duser:="yansapi"}
: ${yansapi_dgroup:="yans"}
: ${yansapi_basepath:="@LOCALSTATEDIR@/yans"}

name="yansapi"
rcvar="yansapi_enable"
command="@BINDIR@/${name}"
command_args="-u ${yansapi_duser} -g ${yansapi_dgroup} -b ${yansapi_basepath}"
command_interpreter="yans"
start_precmd="yansapi_precmd"
pidfile="${yansapi_basepath}/yansapi.pid"

yansapi_precmd()
{
  install -m 775 -o root -g yans -d "${yansapi_basepath}" || {
    err 1 "Failed to create ${yansapi_basepath}"
  }

  if [ -f "$pidfile" ]; then
    local pid=$(check_pidfile "$pidfile" "$command")
    if [ -z "$pid" ]; then
      warn "Removing stale pidfile"
      rm "$pidfile"
    else
      # Shouldn't happen due to the way rc works, but check for it anyway
      err 1 "Already running at $pid"
    fi
  fi

  return 0
}

load_rc_config $name

# FIXME: yansapi is shebanged with #!/usr/bin/env looking for yans, and
#        /usr/local/bin is not a part of PATH at this point
PATH="${PATH}:/usr/local/bin" ; export PATH

run_rc_command "$1"
